# ansible/templates/upstart.conf
description "{{ project_name }}"

start on (filesystem)
stop on runlevel [016]

respawn
console log
setuid nobody
setgid nogroup
chdir {{ application_dir }}/{{ project_name }}

env DJANGO_SETTINGS_MODULE={{ project_name}}.settings.production
exec {{ virtualenv_bin_path}}/gunicorn {{ project_name }}.wsgi:application
setuid www-data
